# templates/default/monit.upstart.erb
# Generated by Chef
description "Monit service manager"

limit core unlimited unlimited

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn

exec <%= @binary %> -c <%= @conf_file %>

pre-stop exec <%= @binary %> -c <%= @conf_file %> quit
