# templates/default/monit.init.erb
#! /bin/sh
#
<% case @platform_family %>
<% when 'rhel' %>
# monit         Monitor Unix systems
#
# chkconfig:    2345 98 02
# description:  Monit is a utility for managing and monitoring processes,
#               files, directories and filesystems on a Unix system. 
# processname:  monit
# pidfile:      /var/run/monit.pid
# config:       <% @conf_file %>
<% when 'debian' %>
### BEGIN INIT INFO
# Provides:          monit
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Should-Start:      $all
# Should-Stop:       $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: service and resource monitoring daemon
# Description:       monit is a utility for managing and monitoring
#                    processes, programs, files, directories and filesystems
#                    on a Unix system. Monit conducts automatic maintenance
#                    and repair and can execute meaningful causal actions
#                    in error situations.
### END INIT INFO
<% end %>

MONIT=<%= @binary %>
CONFIG=<%= @conf_file %>
PID=/var/run/monit.pid

[ -x $MONIT ] || exit 0

RETVAL=0

# See how we were called.
case "$1" in
  start)
        echo -n "Starting monit: "
        $MONIT -c $CONFIG -p $PID
        RETVAL=$?
        echo
        ;;
  stop)
        echo -n "Stopping monit: "
        $MONIT -c $CONFIG -p $PID quit
        RETVAL=$?
        echo
        ;;
  restart)
  	$0 stop
	$0 start
	RETVAL=$?
	;;
  status)
        $MONIT -c $CONFIG -p $PID status
        RETVAL=$?
        ;;
  reload)
        echo -n "Reloading monit: "
        $MONIT -c $CONFIG -p $PID reload
        RETVAL=$?
        echo
        ;;
  *)
	echo "Usage: $0 {start|stop|restart|condrestart|status}"
	exit 1
esac

exit $RETVAL
